# DynamicDevelopment/resources/tornado.8001.conf
[program:tornado8001]
command=/usr/bin/python3 /home/ubuntu/web/web.py -p 8001
directory=/tmp
autostart=true
autorestart=true
startretries=0
stderr_logfile=/var/log/tornado.8001.err.log
stdout_logfile=/var/log/tornado.8001.out.log
user=www-data

