# Dynamic development/resources/tornado.8000.conf
[program:tornado8000]
command=/usr/bin/python3 /home/ubuntu/web/web.py -p 8000
directory=/tmp
autostart=true
autorestart=true
startretries=0
stderr_logfile=/var/log/tornado.8000.err.log
stdout_logfile=/var/log/tornado.8000.out.log
user=www-data

